# crontab
* * * * * echo "Hello world" > /tmp/foo.log
